# files/etc/uci-defaults/nfsd-disable
#!/bin/sh

/etc/init.d/nfsd disable
/etc/init.d/nfsd stop

exit 0
